# docker/common/supervisord.conf
# Baliverne desktop stack (Neko-style): full Xorg with xorg.conf (dummy driver, modelines) → PulseAudio → Openbox → app.
# See https://github.com/m1k1o/neko/blob/master/runtime/supervisord.conf

[supervisord]
nodaemon=true
user=root
loglevel=info

[program:x-server]
environment=DISPLAY="%(ENV_DISPLAY)s"
command=/usr/bin/Xorg %(ENV_DISPLAY)s -config /etc/baliverne/xorg.conf -noreset -nolisten tcp
autorestart=true
priority=200
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:pulseaudio]
environment=DISPLAY="%(ENV_DISPLAY)s",PULSE_SERVER="%(ENV_PULSE_SERVER)s"
command=/bin/sh -c 'sleep 2 && /usr/bin/pulseaudio --log-level=error --disallow-module-loading --disallow-exit --exit-idle-time=-1 -F /etc/baliverne/default.pa'
autorestart=true
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:openbox]
environment=DISPLAY="%(ENV_DISPLAY)s"
command=/bin/sh -c 'sleep 2 && /usr/bin/openbox'
autorestart=true
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:baliverne]
environment=HOME="/app",DISPLAY="%(ENV_DISPLAY)s",PULSE_SERVER="%(ENV_PULSE_SERVER)s",PATH="/usr/local/bin:/usr/bin:/bin"
command=/bin/sh -c 'sleep 3 && exec node /app/index.js'
directory=/app
autorestart=true
priority=800
stopsignal=INT
stopwaitsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
